# fpm-package-builder/utils/contributor/extras/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade|deconfigure)
        echo "Stopping contributoor service..."
        systemctl stop contributoor.service || true # Don't fail if already stopped

        echo "Disabling contributoor service..."
        systemctl disable contributoor.service || true # Don't fail if not enabled
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# Debhelper placeholders below - leave them as is
# Automatically added by dh_systemd_stop/13.3.4
if [ "$1" = remove ] || [ "$1" = upgrade ] || [ "$1" = deconfigure ] ; then
	: # Rely on explicit stop/disable above
fi
# End automatically added section

exit 0
