# fpm-package-builder/l1-clients/execution-layer/nimbus-ec/extras/nimbus-execution-holesky.service
[Unit]
Description=Nimbus Execution client daemon (Holesky)
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/ethereum/nimbus-execution-holesky.conf
ExecStart=/usr/bin/nimbus_execution_client $ARGS
Restart=on-failure
RestartSec=5
User=ethereum
SyslogIdentifier=nimbus-execution-holesky
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target