# fpm-package-builder/teku/extras/teku-besu.service
[Unit]
Description=Ethereum Teku client daemon
After=network.target

[Service]
EnvironmentFile=/etc/ethereum/kiln/teku-besu.conf
ExecStart=/usr/bin/teku-kl $ARGS
Restart=always
User=ethereum

[Install]
WantedBy=multi-user.target