# fpm-package-builder/l2-clients/linea/extras/linea-advanced.service
[Unit]
Description=Linea Managed Besu Execution Client (Advanced Profile)
After=network.target

[Service]
EnvironmentFile=/etc/ethereum/linea-advanced.conf
ExecStart=/usr/bin/linea-besu $ARGS
Restart=always
RestartSec=5
User=ethereum
Type=simple

[Install]
WantedBy=multi-user.target
