# fpm-package-builder/l1-clients/consensus-layer/prysm/extras/prysm-validator-hoodi.service
[Unit]
Description=Prysm Validator daemon (Hoodi)

After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/ethereum/prysm-validator-hoodi.conf
ExecStart=/usr/bin/validator $ARGS

Restart=on-failure
RestartSec=5
User=ethereum
SyslogIdentifier=prysm-validator-hoodi

LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target

