# fpm-package-builder/l1-clients/consensus-layer/grandine/extras/grandine-validator-hoodi.service
[Unit]
Description=Grandine Beacon chain daemon
After=network.target

[Service]
EnvironmentFile=/etc/ethereum/grandine-validator-hoodi.conf
ExecStart=/usr/bin/grandine $ARGS
Restart=always
User=ethereum
SyslogIdentifier=grandine-validator-hoodi

[Install]
WantedBy=multi-user.target

