# fpm-package-builder/infra/dvt/obol/extras/prysm-validator-obol.service
[Unit]
Description=Prysm Validator Client (Obol)
After=network-online.target
Wants=network-online.target
Requires=prysm-beacon.service

[Service]
EnvironmentFile=/etc/ethereum/dvt/prysm-validator-obol.conf
ExecStart=/usr/bin/validator $ARGS
Restart=always
User=ethereum
SyslogIdentifier=prysm-validator-obol

[Install]
WantedBy=multi-user.target

