# fpm-package-builder/infra/dvt/obol/extras/teku-validator-hoodi-obol.service
[Unit]
Description=Teku Validator Client (Obol Hoodi)
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/ethereum/dvt/teku-validator-hoodi-obol.conf
ExecStart=/usr/bin/teku $ARGS
Restart=always
User=ethereum
SyslogIdentifier=teku-validator-hoodi-obol

[Install]
WantedBy=multi-user.target
