# fpm-package-builder/utils/ethereumonarm-utils/sources/usr/lib/systemd/system/ethereum-swappiness-enforce.service
[Unit]
Description=Enforce Ethereum on ARM sysctl settings
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/ethereum-swappiness-enforce
