# fpm-package-builder/l1-clients/execution-layer/nethermind/extras/nethermind.service
[Unit]
Description=Nethermind Ethereum Daemon
After=network.target

[Service]
User=ethereum
EnvironmentFile=/etc/ethereum/nethermind.conf
ExecStart=/usr/bin/nethermind $ARGS
Restart=on-failure
WorkingDirectory=/home/ethereum
LimitNOFILE=1000000
SyslogIdentifier=nethermind

[Install]
WantedBy=default.target
