# fpm-package-builder/l1-clients/execution-layer/nethermind/extras/nethermind-base.service
[Unit]
Description=Nethermind Ethereum ECL Daemon (Base)
After=network-online.target
Wants=network-online.target

[Service]
User=ethereum
EnvironmentFile=/etc/ethereum/nethermind-base.conf
ExecStart=/usr/bin/nethermind $ARGS
Restart=on-failure
RestartSec=5
WorkingDirectory=/home/ethereum
LimitNOFILE=1000000
LimitNPROC=65536
SyslogIdentifier=nethermind-base

[Install]
WantedBy=multi-user.target