# fpm-package-builder/l1-clients/execution-layer/erigon3/extras/erigon.service
[Unit]
Description=Ethereum Erigon client daemon
After=network.target

[Service]
Environment="GODEBUG=netdns=go"
EnvironmentFile=/etc/ethereum/erigon.conf
ExecStart=/usr/bin/erigon $ARGS
Restart=always
User=ethereum

[Install]
WantedBy=multi-user.target

