# fpm-package-builder/l1-clients/consensus-layer/nimbus/extras/nimbus-beacon.service
[Unit]
Description=Nimbus Beacon chain daemon (Mainnet)
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/ethereum/nimbus-beacon.conf
ExecStart=/usr/bin/nimbus_beacon_node $ARGS
Restart=on-failure
RestartSec=5
User=ethereum
SyslogIdentifier=nimbus-beacon
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target