# fpm-package-builder/utils/ethereumonarm-armbian-extras/sources/etc/cron.d/armbian-truncate-logs
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/5 * * * * root /usr/lib/armbian/armbian-truncate-logs
