# fpm-package-builder/l1-clients/consensus-layer/lodestar/extras/lodestar-beacon-hoodi.service
[Unit]
Description=Lodestar Beacon chain daemon (Hoodi)
# Waits for a full network connection
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/ethereum/lodestar-beacon-hoodi.conf
# Fix: Includes $NODE_OPTIONS for memory management
ExecStart=/usr/bin/lodestar $NODE_OPTIONS $ARGS
# Reliable restart policy
Restart=on-failure
RestartSec=5
User=ethereum
SyslogIdentifier=lodestar-beacon-hoodi
# Increased resource limits
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target