# fpm-package-builder/utils/ethereumonarm-config-sync/extras/ethereumonarm-config-sync.timer
[Unit]
Description=Runs ethereumonarm-config-sync every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target