# fpm-package-builder/infra/liquid-staking/lido/extras/teku-validator-hoodi-lido.service
[Unit]
Description=Teku Validator Client (Lido Hoodi)
After=network.target

[Service]
EnvironmentFile=/etc/ethereum/ls/teku-validator-hoodi-lido.conf
ExecStart=/usr/bin/teku $ARGS
Restart=always
User=ethereum

[Install]
WantedBy=multi-user.target
