# fpm-package-builder/prysm/extras/prysm-beacon-prater.service
[Unit]
Description=Prysm Beacon chain daemon
After=network.target

[Service]
EnvironmentFile=/etc/ethereum/prysm-beacon-prater.conf
ExecStart=/usr/bin/beacon-chain $ARGS
Restart=always
User=ethereum

[Install]
WantedBy=multi-user.target

