# fpm-package-builder/utils/ethereumonarm-config-sync/sources/etc/cron.d/ethereumonarm-config-sync
# /etc/cron.d/ethereumonarm-config-sync: crontab entries for the Ethereum on ARM sync
#
# This script is executed every hour.
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 * * * * root test -x /usr/local/sbin/ethereumonarm-config-sync.sh && /usr/local/sbin/ethereumonarm-config-sync.sh